# deploy/systemd/aeron-driver.service
[Unit]
Description=Aeron Media Driver (IPC shared-memory transport)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=arbx
ExecStartPre=/usr/bin/podman pull localhost/aeron-driver:latest
ExecStart=/usr/bin/podman run --rm --name aeron-driver \
    --ipc=host \
    --ulimit memlock=-1:-1 \
    -v /dev/shm:/dev/shm \
    localhost/aeron-driver:latest
ExecStop=/usr/bin/podman stop aeron-driver
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
