# deploy/systemd/arbx-sidecar-fubon.service
[Unit]
Description=ArbX Sidecar - Fubon (stream 1002)
After=network-online.target aeron-driver.service
Wants=network-online.target
Requires=aeron-driver.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=arbx
EnvironmentFile=/etc/arbx/.env
ExecStartPre=/usr/bin/podman pull localhost/sidecar-fubon:latest
ExecStart=/usr/bin/podman run --rm --name arbx-sidecar-fubon \
    --network=host \
    --ipc=host \
    --ulimit memlock=-1:-1 \
    --env-file /etc/arbx/.env \
    -v /dev/shm:/dev/shm \
    -v /etc/arbx/config:/etc/arbx/config:ro \
    --mount type=tmpfs,destination=/certs,tmpfs-size=1M \
    -v /etc/arbx/certs/fubon:/certs:ro \
    localhost/sidecar-fubon:latest \
    /etc/arbx/config/fubon.yaml
ExecStop=/usr/bin/podman stop arbx-sidecar-fubon
Restart=on-failure
RestartSec=5
# Resource caps (review §3.10).
MemoryMax=2G
LimitNOFILE=65536
LimitMEMLOCK=infinity
# Sandbox hardening.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=yes
LockPersonality=yes
RestrictSUIDSGID=yes
ReadWritePaths=/dev/shm
ReadOnlyPaths=/etc/arbx

[Install]
WantedBy=multi-user.target
